# meta-wiwynn/meta-sv310a/recipes-phosphor/sv310a-boot-option/boot-option/boot-option.service
[Unit]
Description=Boot option
Before=phosphor-ipmi-host.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/boot-option

[Install]
WantedBy=multi-user.target
