# meta-wiwynn/meta-sv310a/recipes-phosphor/sv310a-leaky-bucket/leaky-bucket/service/xyz.openbmc_project.leaky-bucket.service
[Unit]
Description=Leaky Bucket Service
StopWhenUnneeded=false
Before=xyz.openbmc_project.Logging.IPMI.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/leaky-bucket
SyslogIdentifier=leaky-bucket

[Install]
WantedBy=multi-user.target
