# meta-wiwynn/meta-sv310a/recipes-phosphor/sv310a-service-oem/files/sv310a-bmcenv-init.service
[Unit]
Description=SV310A BMC Environment Initialization

[Service]
Type=oneshot
ExecStart=/usr/sbin/sv310a-bmcenv-init.sh
SyslogIdentifier=sv310a-bmcenv-init

[Install]
WantedBy=sysinit.target
