# meta-wiwynn/meta-sv310a/recipes-phosphor/ipmi/phosphor-ipmi-kcs/phosphor-ipmi-kcs@.service.in
[Unit]
Description=Phosphor IPMI KCS DBus Bridge
After=multi-user.target

[Service]
Type=notify
Restart=always
ExecStart=@BIN@ -c "%i"
SyslogIdentifier=kcsbridged-%i

[Install]
WantedBy=multi-user.target
RequiredBy=
