# meta-wiwynn/meta-sv310a/recipes-phosphor/gpio/sv310a-gpio-monitor-register/SetNextNodeAbsent.service
[Unit]
Description=Set Next Node Absent

[Service]
Type=oneshot

ExecStart=echo "Another node absent detected"

ExecStopPost=/usr/bin/busctl set-property xyz.openbmc_project.GpioMonitor /xyz/openbmc_project/gpio/event xyz.openbmc_project.GpioEvent NEXT_NODE_P3V3_STBY_PG s "Idle"

