# meta-wiwynn/meta-sv310a/recipes-phosphor/ipmi/phosphor-ipmi-flash/update-image.service
[Unit]
Description=Update BMC Image File

[Service]
Type=oneshot
ExecStart=/usr/bin/env update_image.sh

[Install]
RequiredBy=phosphor-ipmi-flash-bmc-update.target

