# meta-wiwynn/meta-sv310a/recipes-extended/rsyslog/rsyslog/rotate-event-logs.service
[Unit]
Description=Rotates the event logs

[Service]
Type=simple
ExecStart=/usr/bin/rotate-event-logs

[Install]
WantedBy=multi-user.target
