# meta-wiwynn/meta-sv310a/recipes-phosphor/sv310a-gpio-init/files/sv310a-gpio-init.service
[Unit]
Description=SV310A GPIO Initialization

[Service]
Type=oneshot
ExecStart=/usr/sbin/sv310a-gpio-init.sh
SyslogIdentifier=sv310a-gpio-init

[Install]
WantedBy=sysinit.target
